# 03_zookeeper/files/supervisor-zookeeper.conf
[program:zookeeper]
user=zookeeper
command=/usr/lib/zookeeper/bin/zkServer.sh start-foreground
autorestart=true
startsecs=45
startretries=3
